# deploy/systemd/allocator.service
[Unit]
Description=Portfolio Allocator (Fleet Risk Monitor)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/allocator
ExecStart=/bin/bash -c '/usr/bin/python3 -u /home/ubuntu/allocator/portfolio_allocator.py >> /home/ubuntu/allocator/allocator.log 2>&1'
Restart=always
RestartSec=10
MemoryMax=200M

[Install]
WantedBy=multi-user.target
