# deploy/systemd/gap-bot.service
[Unit]
Description=Gap Fade/Continuation Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bot_gap
ExecStart=/bin/bash -c "/usr/bin/python3 -u /home/ubuntu/bot_gap/gap_bot.py >> /home/ubuntu/bot_gap/gap.log 2>&1"
Restart=always
RestartSec=10
MemoryMax=200M

[Install]
WantedBy=multi-user.target
